# deploy/systemd/assessmentportal-worker.service
[Unit]
Description=Assessment Portal Assessment Worker
After=network-online.target assessmentportal-gunicorn.service
Wants=network-online.target

[Service]
Type=simple
User=assessmentportal
Group=assessmentportal
WorkingDirectory=/opt/assessmentportal/backend
EnvironmentFile=/etc/assessmentportal/assessmentportal.env
Environment=HOME=/var/lib/assessmentportal
ExecStart=/opt/assessmentportal/.venv/bin/python manage.py run_queued_assessments
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/assessmentportal

[Install]
WantedBy=multi-user.target
